# iptv-server.service
[Unit]
Description=IPTV Rebroadcast Server
After=network.target

[Service]
Type=simple
User=CURRENT_USER
WorkingDirectory=/path/to/project
ExecStart=/usr/bin/node /path/to/project/server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=iptv-server

# Environment file (create from .env.example)
EnvironmentFile=-/path/to/project/.env

# Default environment (override via .env file)
Environment=NODE_ENV=production
Environment=PORT=3000

# Resource limits
LimitNOFILE=65536
LimitNPROC=512

[Install]
WantedBy=multi-user.target
